# srcs/influxdb/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command = sh -c "influxd && kill -s SIGTERM $(cat supervisord.pid)"

[program:telegraf]
command = sh -c "/usr/bin/telegraf -config /etc/telegraf/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"
